# etc/runit/rc/rc.local
#!/bin/bash
#
# /etc/runit/rc/rc.local: Local multi-user start-up script.
#


